# autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
Comment=Simple X hotkey daemon map F1 to open alacritty
Exec=/usr/bin/sxhkd
Terminal=false
Type=Application
