# autostart/xcape.desktop
[Desktop Entry]
Type=Application
Name=xcape
Exec=setxkbmap -option ctrl:nocaps && xcape -e 'Control_L=Escape'
